# deploy/systemd/chrome-jfn.service
[Unit]
Description=Chrome headless CDP bridge (porta 9222) para coleta TFE/SIAFE do JFN
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStartPre=/bin/mkdir -p /tmp/chrome-jfn
ExecStart=/usr/bin/google-chrome --headless=new --no-sandbox --disable-gpu --no-first-run --remote-debugging-address=127.0.0.1 --remote-debugging-port=9222 --user-data-dir=/tmp/chrome-jfn about:blank
Restart=always
RestartSec=5
# baixa prioridade: ponte só é usada sob demanda pelo JFN
Nice=10

[Install]
WantedBy=default.target
