# deploy/systemd/jfn-ronda.service
[Unit]
Description=JFN — ronda (checa serviço + alertas e notifica)
After=jfn.service

[Service]
Type=oneshot
WorkingDirectory=%h/JFN
Environment=JFN_DATA_DIR=%h/JFN/data
Environment=JFN_BASE_URL=http://127.0.0.1:8000
Environment=JFN_SERVICE=jfn.service
ExecStart=%h/JFN/.venv/bin/python %h/JFN/tools/ronda.py
