# deploy/systemd/jfn-sancoes.timer
[Unit]
Description=JFN — CEIS/CNEP semanal (domingo 05:40, antes do ciclo do Núcleo)

[Timer]
OnCalendar=Sun 05:40
Persistent=true

[Install]
WantedBy=timers.target
