# deploy/systemd/jfn-ronda.timer
[Unit]
Description=JFN — agenda a ronda a cada 10 minutos

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
AccuracySec=30s
Persistent=true
Unit=jfn-ronda.service

[Install]
WantedBy=timers.target
