# deploy/systemd/massare-market.service
[Unit]
Description=Massare — cenarios de pregao multi-horizonte (ao vivo)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/JFN
ExecStart=%h/JFN/.venv/bin/python -m massare.market
Nice=10
