# _SANDBOX/gcp/propostas/2026-06-04-yoda-deploy/yoda.service
[Unit]
Description=Yoda - Hermes Telegram Gateway (projeto Hermes)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jfelippebethlem
Group=jfelippebethlem
WorkingDirectory=/home/jfelippebethlem/hermes-agent
Environment=HERMES_HOME=/home/jfelippebethlem/.hermes
Environment=PYTHONIOENCODING=utf-8
ExecStart=/home/jfelippebethlem/hermes-agent/venv/bin/python -m hermes_cli.main gateway run
Restart=always
RestartSec=10
KillSignal=SIGTERM
TimeoutStopSec=210

[Install]
WantedBy=multi-user.target
