# deploy/jfn-fisc-refresh.service
[Unit]
Description=JFN — refresh semanal da fiscalização (emendas RJ + PIX + favorecidos + PNCP municipal Rio)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/JFN
ExecStart=%h/JFN/.venv/bin/python tools/fisc_coletar_tudo.py --incremental --pausa 1.0
TimeoutStartSec=14400
Nice=10
